# server/etc/rc.local
# Start OpenVPN on boot
/usr/local/sbin/openvpn --daemon --config /etc/openvpn/server.conf